# deployment/fastapi-fileserver.service
[Unit]
Description=FastAPI File Server API
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/fastapi-fileserver
EnvironmentFile=/etc/fastapi-fileserver.env
ExecStart=/var/www/fastapi-fileserver/venv/bin/gunicorn main:app -w 4 -k uvicorn.workers.UvicornWorker --bind 127.0.0.1:8000 --keep-alive 5 --max-requests 1000 --max-requests-jitter 100
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target